# systemd/coco-agent.service
[Unit]
Description=Coco Native Agent (manual runs only)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/coco-device
ExecStart=/usr/local/bin/coco-native-agent-boot.sh
Restart=on-failure
RestartSec=10
StandardOutput=append:/var/log/coco/agent.log
StandardError=append:/var/log/coco/agent.log

[Install]
WantedBy=multi-user.target
