# systemd/coco-agent-scheduler.service
[Unit]
Description=Coco scheduled session runner (twice daily)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=pi
WorkingDirectory=/home/pi/coco-device
ExecStart=/usr/local/bin/run-scheduled-session.sh
Environment=COCO_RUN_USER=pi

[Install]
WantedBy=multi-user.target
